# ansible/roles/kubernetes-common/templates/etc/systemd/system/kubelet.service.d/09-vagrant.conf
[Service]
Environment="KUBELET_EXTRA_ARGS= --node-ip={{ kubernetes_node_ip }}"
