# ansible/roles/aws/templates/etc/systemd/system/kubelet.service.d/09-cloud-provider.conf
[Service]
Environment="KUBELET_CLOUD_PROVIDER=--cloud-provider=aws"
