# ansible/roles/kubernetes/templates/etc/systemd/system/kubelet.service.d/10-extra-args.conf
[Service]
Environment="KUBELET_EXTRA_ARGS= $KUBELET_CLOUD_PROVIDER"
